# root/etc/s6-overlay/s6-rc.d/svc-wireguard-launcher/run
#!/usr/bin/with-contenv bash
exec java -jar /lib/wireguard-launcher.jar
